# app/crontab
# Run weather_update.py every Monday at 5 am
0 5 * * 1 python /app/weather_update.py

# Run main.py every day at 5 am
0 5 * * * python /app/main.py


